# deploy/unit-files/gaia.timer
# gaia.timer
[Unit]
Description=Gaia Hub timer

[Timer]
OnStartupSec=60

[Install]
WantedBy=timers.target